# systemd/cluster_guard.timer
[Unit]
Description=Run Cluster Guard every minutes

[Timer]
OnCalendar=*:0/2

[Install]
WantedBy=timers.target
